# tilde/dot_config/systemd/user/exact_waybar.service.d/requires_mpris.conf
[Unit]
Requires=mpd-mpris.service
After=mpd-mpris.service
